# supervisor/http-command.conf
[program:http-command]
command = /home/pi/.virtualenvs/misc/bin/python /home/pi/.virtualenvs/misc/src/http-reboot.py

autostart = true
autorestart = true

