# supervisor/minidlna-supervisor.conf
[program:minidlna]
directory = /usr/local/sbin
command = /usr/local/sbin/minidlnad -R -f /home/odroid/.minidlna/minidlna.conf

autostart = true
#autorestart = true

