# supervisor/reboot-pi.conf
[program:reboot-pi]
command = sh /home/odroid/reboot-all-pi.sh 

autostart = true
autorestart = false
startsecs = 1

