# supervisor/trim_old_files.conf
[program:trim_old_files]
directory = /home/pi/.virtualenvs/misc/src
command = /home/pi/.virtualenvs/misc/bin/python trim-oldfiles.py

autostart = true
autorestart = true

