# supervisor/weather.conf
[program:weather]
command = /home/pi/.virtualenvs/misc/bin/python /home/pi/.virtualenvs/misc/src/weather.py

autostart = true
autorestart = true

